# configs/supervisor.conf.d/crond.conf
[program:cron]
command=/usr/sbin/crond -f
autorestart=true
